# supervisor_cron.conf
[program:cron]
command=/root/runcron.sh
user=root
autostart=true
autorestart=true
redirect_stderr=true
