# supervisor_mssql.conf
[program:mssql]
command=/opt/mssql/bin/sqlservr
user=root
autostart=true
autorestart=true
redirect_stderr=true
